# ops/chainya-backup.service
[Unit]
Description=Backup Chainya orders database

[Service]
Type=oneshot
User=root
UMask=0077
ExecStart=/opt/chainya-shop/.venv/bin/python /opt/chainya-shop/ops/backup-chainya.py
Nice=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/chainya-shop /var/lib/chainya-shop
